# macos/com.user.power-app-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.power-app-manager</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/hb/dotfiles/exec/power-app-manager</string>
        <string>monitor</string>
    </array>

    <key>StartInterval</key>
    <integer>30</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/hb/Library/Logs/power-app-manager-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/hb/Library/Logs/power-app-manager-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
